# launchd/com.trader-obsidian.weekly-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trader-obsidian.weekly-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/al/Documents/ClaudeCode/obsidiantrader/scripts/weekly_review.py</string>
        <string>--notify</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/al/Documents/ClaudeCode/obsidiantrader</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/al/Documents/ClaudeCode/obsidiantrader</string>
    </dict>

    <!-- 每周六上午 10:00 运行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/trader-weekly-review.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/trader-weekly-review.err</string>
</dict>
</plist>
