# launchd/com.trader-obsidian.inbox.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trader-obsidian.inbox</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/al/Documents/Alma/trader/trader-obsidian/scripts/scan_inbox.py</string>
        <string>--notify</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/al/Documents/Alma/trader/trader-obsidian</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/al/Documents/Alma/trader/trader-obsidian</string>
    </dict>
    
    <!-- 每30分钟运行一次 -->
    <key>StartInterval</key>
    <integer>1800</integer>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/trader-inbox.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/trader-inbox.err</string>
</dict>
</plist>
